# nrbdaq.service
[Unit]
After=network.target

[Service]
ExecStart=/home/gaw/git/nrbdaq/nrbdaq.sh

[Install]
WantedBy=default.target